# docs/pi_config/bitos-device.service
[Unit]
Description=BITOS Device Runtime
After=bitos-server.service
Requires=bitos-server.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bitos
EnvironmentFile=/etc/bitos/secrets
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_FBDEV=/dev/fb0
Environment=BITOS_DISPLAY=st7789
Environment=BITOS_AUDIO=hw:0
Environment=BITOS_BUTTON=gpio
Environment=SERVER_URL=http://localhost:8000
Environment=WHISPLAY_DRIVER_PATH=/home/pi/Whisplay/Driver
Environment=PYTHONPATH=/home/pi/bitos
Environment=PYTHONWARNINGS=ignore::RuntimeWarning
ExecStartPre=/bin/bash -c "for i in $(seq 1 10); do curl -sf http://localhost:8000/health && exit 0; sleep 1; done; echo 'Server not ready after 10s, starting device anyway'; exit 0"
ExecStart=/home/pi/bitos/.venv/bin/python -m device.main
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource limits
MemoryMax=256M
MemorySwapMax=128M
CPUQuota=80%
OOMScoreAdjust=500

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/bitos/device/data /home/pi/bitos/logs /tmp
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
